# scripts/packaging/mavkit/debian/mavkit-baker.mavkit-accuser-next.service
[Unit]
Description=Mavkit accuser
Documentation=https://protocol.mavryk.org/
After=network.target

[Service]
EnvironmentFile=-/etc/default/mavkit-baker-next
EnvironmentFile=-/etc/default/mavkit-accuser-next
WorkingDirectory=/var/tezos
Restart=on-failure
User=tezos
ExecStart=/bin/sh -c "/usr/bin/mavkit-accuser-$PROTOCOL run $RUNTIME_OPTS"
StandardOutput=append:/var/log/tezos/accuser-next.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
